# squashdir
#!/sbin/openrc-run
description="Squash the upperdirectory on shutdown"
rc_shutdown=yes

depend() {
	need local
	need localmount
}

start() {
	ebegin "Squashing Upperdir..."
    /sbin/squash-upperdir	
	eend $?
}
